# scripts/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# connect to lan
wpa_supplicant -Dwired -ieth0 -c/etc/wpa_supplicant.conf -B &
dhcpcd eth0 &

# enable direct monitor brightness adjustment
chmod a+rw /sys/class/backlight/acpi_video0/brightness &

# bluetooth is dead
rfkill block bluetooth &

